# launchd/com.adames.workspace.autohide.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adames.workspace.autohide</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adames/.local/bin/ws-autohide</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>LSUIElement</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/Users/adames/.cache/workspace/ws-autohide.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/adames/.cache/workspace/ws-autohide.err</string>

    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/adames</string>
        <key>PATH</key>
        <string>/Users/adames/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
